# systemd/livebox-reboot.service
[Unit]
Description=Reboot Livebox

[Service]
Type=oneshot
ExecStart=/usr/local/bin/livebox-reboot
